# vendor/stash.service
[Unit]
Description=Wayland clipboard manager with fast persistent history and multi-media support
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target

[Service]
Type=simple
ExecStart=stash watch
Restart=on-failure

[Install]
WantedBy=graphical-session.target
